# systemd/floppa-server.service
[Unit]
Description=Floppa VPN Server (Bot + Admin API)
After=network-online.target postgresql.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=floppa
WorkingDirectory=/opt/floppa-vpn
ExecStart=/opt/floppa-vpn/bin/floppa-server
Environment=FLOPPA_CONFIG=/etc/floppa-vpn/config.toml
Environment=FLOPPA_SECRETS=/etc/floppa-vpn/secrets.toml
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
